# themes/luci-theme-wifimedia/root/etc/uci-defaults/30_luci-theme-netgear
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Wifimedia=/luci-static/wifimedia
	set luci.main.mediaurlbase=/luci-static/wifimedia
	set luci.main.lang=en
	commit luci
EOF
exit 0
